# recipes-core/fd-dsp_0.1/files/fd-dsp.service
[Unit]
Description=Fiberdyne DSP

[Service]
Type=forking
ExecStart=/usr/bin/fd-dsp

[Install]
WantedBy=default.target
