# recipes-core/fd-dsp_0.1/files/fd-dsp-pulse.service
[Unit]
Description=Fiberdyne DSP PulseAudio Plugin
After=fd-dsp.service

[Service]
Type=oneshot
ExecStart=/usr/bin/fd-dsp-pulse

[Install]
WantedBy=default.target
